# scripts/launchd/com.rp.wochenplaene.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.rp.wochenplaene.update</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/rp/Documents/projekte/ernaehrung/automation/wochenplaene/.venv/bin/python</string>
      <string>/Users/rp/Documents/projekte/ernaehrung/automation/wochenplaene/scripts/update_weekly_plan.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/rp/Documents/projekte/ernaehrung/automation/wochenplaene</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/rp/Documents/projekte/ernaehrung/automation/wochenplaene/logs/update_weekly_plan.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rp/Documents/projekte/ernaehrung/automation/wochenplaene/logs/update_weekly_plan.err</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
  <key>Hour</key>
  <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </dict>
</plist>
